# updateFreeDNS.service
[Unit]
Description=Update FreeDNS periodically

[Service]
ExecStart=/usr/bin/updateFreeDNS

[Install]
WantedBy=multi-user.target 
